# redis/rootfs/etc/s6-overlay/s6-rc.d/init-redis/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Redis
# Redis init script, runs before any other service
# ==============================================================================

bashio::log.info "Configuring Redis..."

port="$(bashio::addon.port 6379)"
if ! bashio::var.has_value "${port}"; then
  port=6379
fi


bashio::var.json \
  port "^$port" \
  ssl "^$(bashio::config 'ssl')" \
  certfile "$(bashio::config 'certfile')" \
  keyfile "$(bashio::config 'keyfile')" \
  db_password "$(bashio::config 'db_password')" \
  | tempio \
      -template /etc/redis/config.gtpl \
      -out /config/config.conf

if [ "$(bashio::config 'log_level')" == "debug" ]; then
  bashio::log.info "$(cat /config/config.conf)"
fi
