# redis/rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: NRedis8N
# Runs Redis script
# ==============================================================================

bashio::log.info "Starting Redis..."

redis-server